# bin/decrypt.service
[Unit]
Description=Decryption

[Service]
ExecStart=/home/vennila/Documents/cryptography/bin/decrypt.sh

[Install]
WantedBy=multi-user.target
